# dist/bluephysics/share/gvfs/mounts/ftpis.mount
[Mount]
Type=ftpis
Exec=/usr/lib/gvfs/gvfsd-ftp
AutoMount=false
Scheme=ftpis
DefaultPort=990
HostnameIsInetAddress=true
