# dist/bluephysics/share/gvfs/mounts/nfs.mount
[Mount]
Type=nfs
Exec=/usr/lib/gvfs/gvfsd-nfs
AutoMount=false
Scheme=nfs
